# supervisord.conf
[supervisord]
nodaemon=true

[program:phddns]
command=phddns start
autostart=true


[program:status]
command=phddns status
