# supervisord.conf
[supervisord]
nodaemon=true

[program:phddns]
command=phddns start
autostart=true
autorestart=true


[program:status]
command=phddns status
autostart=true
autorestart=false